# smedjan/scripts/health-dashboard.service
[Unit]
Description=Smedjan factory health HTML dashboard generator
Documentation=https://github.com/anders/agentindex/blob/main/smedjan/scripts/health_dashboard.py
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=smedjan
Group=smedjan
WorkingDirectory=/home/smedjan/agentindex
Environment=PYTHONPATH=/home/smedjan/agentindex
Environment=SMEDJAN_CONFIG_DIR=/home/smedjan/smedjan/config
ExecStart=/usr/bin/python3 -m smedjan.scripts.health_dashboard
StandardOutput=append:/home/smedjan/smedjan/worker-logs/health-dashboard.log
StandardError=inherit
TimeoutStartSec=2m
Nice=15
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/smedjan/smedjan
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
